# web_bbs.conf
[program:web_bbs]
command=/usr/local/bin/gunicorn wsgi --bind 0.0.0.0:2000 --workers 3 --worker-class gevent
directory=/var/www/web_bbs
autostart=true
autorestart=true

stdout_logfile=/var/log/supervisor/web_bbs_stdout.log
stderr_logfile=/var/log/supervisor/web_bbs_stderr.log